# configs/xray.rcd
#!/bin/sh
# PROVIDE: xray
# REQUIRE: DAEMON NETWORKING

. /etc/rc.subr

export XRAY_LOCATION_ASSET=/usr/local/etc/xray

name=xray
rcvar=${name}_enable
command="/usr/local/bin/${name}"
pidfile="/var/run/${name}.pid"
start_cmd="${name}_start"
xray_start() {
  /usr/sbin/daemon -cf -p ${pidfile} ${command} run -config /usr/local/etc/xray/config.json
}
load_rc_config $name
run_rc_command $1