# configs/xray.openrc
#!/sbin/openrc-run

export XRAY_LOCATION_ASSET=/usr/local/etc/xray

command="/usr/local/bin/xray"
command_args="run -config /usr/local/etc/xray/config.json"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after sshd firewall
}